# scripts/pi-zero-2w/openrc
#!/sbin/openrc-run

description="Rematrix-gallery Server"
name="rematrix-gallery"

command="uvicorn rematrix_api:app"
command_args="--host 0.0.0.0 --port 80"
directory="/srv/rematrix-gallery"
output_log="/srv/rematrix-gallery/logs/run.log"
error_log="/srv/rematrix-gallery/logs/errors.log"
command_background=true

pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    after hostapd
}